# .config/systemd/user/battery-alert.service
[Unit]
Description=Check battery level and send low battery notification

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/battery-alert
