# .config/systemd/user/rfkill-css-sync.service
[Unit]
Description=Sync swaync radio button highlights on rfkill events
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=%h/.config/hypr/scripts/rfkill-watch.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
